# recipes-rust/simpleserver-rust/files/simpleserver.service
[Unit] 
Description=Simple server in Rust 
 
[Service] 
Type=forking 
ExecStart=/usr/bin/simpleserver-rust 
 
[Install] 
WantedBy=multi-user.target

